# crontab
00 04 * * * root /usr/src/app/update_modules.sh >> /var/log/cron.log 2>&1
# Don't remove the empty line at the end of this file. It is required to run the cron job